# proc_monitor.service
[Unit]
Description=SYSMON kernel module

[Service]
Type=simple
ExecStart=/sbin/insmod /lib/modules/proc_monitor.ko
ExecStop=/sbin/rmmod proc_monitor
RemainAfterExit=yes
LimitCORE=infinity

[Install]
WantedBy=multi-user.target